# ops/puremind-eval.timer
[Unit]
Description=Run pureMind evaluation harness weekly (Saturday 04:00 UTC)

[Timer]
OnCalendar=Sat *-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
