# ops/puremind-metrics.service
[Unit]
Description=pureMind metrics collector -- system health to PostgreSQL
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=puretensorai
Environment=HOME=/home/puretensorai
Environment=PATH=/home/puretensorai/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /home/puretensorai/pureMind/tools/metrics_collector.py
TimeoutStartSec=60
